# com.abrahamadsk.fpt-ami.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.abrahamadsk.fpt-ami</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/abraham/Claude_projects/fpt-mcp/.venv/bin/python3</string>
        <string>-m</string>
        <string>fpt_mcp.ami.handler</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/abraham/Claude_projects/fpt-mcp/src</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/fpt-ami.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/fpt-ami.err</string>
</dict>
</plist>
