# com.abrahamadsk.fpt-mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.abrahamadsk.fpt-mcp</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>-m</string>
        <string>fpt_mcp.server</string>
        <string>--http</string>
        <string>--port</string>
        <string>8090</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/abraham/Claude_projects/fpt-mcp/src</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>SHOTGRID_URL</key>
        <string>https://ableviadsk.shotgrid.autodesk.com</string>
        <key>SHOTGRID_SCRIPT_NAME</key>
        <string>mcp_server</string>
        <key>SHOTGRID_SCRIPT_KEY</key>
        <string>3voplsoawvokdca(mqhhdsXnp</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/fpt-mcp.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/fpt-mcp.err</string>
</dict>
</plist>
